# thenexus.service
[Unit]
Description=TheNexus Dashboard
After=network.target

[Service]
Type=simple
User=azureuser
WorkingDirectory=/home/azureuser/dev/TheNexus
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=thenexus

# Environment
Environment="NODE_ENV=production"
Environment="PATH=/home/azureuser/.linuxbrew/bin:/home/azureuser/.npm-global/bin:/usr/bin:/bin"

[Install]
WantedBy=multi-user.target
